# installers/linux/debian/xai-node.service
[Unit]
Description=XAI Blockchain Node
Documentation=https://docs.xai-blockchain.io
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=xai
Group=xai

# Environment
Environment="XAI_DATA_DIR=/var/lib/xai"
Environment="XAI_CONFIG_DIR=/etc/xai"
Environment="XAI_LOG_DIR=/var/log/xai"
Environment="XAI_NETWORK=testnet"
Environment="PYTHONUNBUFFERED=1"

# Execution
ExecStart=/usr/bin/xai-node --network testnet
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/xai /var/log/xai

# Resource limits
LimitNOFILE=65536
MemoryMax=4G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xai-node

[Install]
WantedBy=multi-user.target
